# dot_config/systemd/user/clipit.service
[Unit]
Description=clipit clipboard manager

[Service]
ExecStart=/usr/bin/clipit
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
