# dot_config/systemd/user/git-wallpaper-clone.service
[Unit]
Description=Update wallpapers from git repository

[Service]
WorkingDirectory=/home/shw/Documents/images/wallpapers
ExecStart=/usr/bin/git pull
Type=oneshot

