# dot_config/systemd/user/wallpaper-rotater.service
[Unit]
Description=Regular wallpaper rotation
RefuseManualStart=yes

[Service]
ExecStart=/home/shw/.local/bin/set-wallpaper.sh
Type=oneshot
